# serve-web/mempalace-mcp-bridge.service
[Unit]
Description=MemPalace MCP HTTP Bridge
After=network.target

[Service]
Type=simple
ExecStart=%h/.venvs/mempalace-copilot/bin/python %h/.config/Code/User/serve-web/mempalace_mcp_http_bridge.py --host 127.0.0.1 --port 3940 --palace %h/.local/share/mempalace/palace
Restart=always
RestartSec=5
WorkingDirectory=%h
Environment=HOME=%h

[Install]
WantedBy=default.target